# deploy/capsulestory.service
[Unit]
Description=Capsule Story API
After=network.target postgresql.service

[Service]
User=capsule
Group=capsule
WorkingDirectory=/home/capsule/capsulestory
Environment="PATH=/home/capsule/capsulestory/venv/bin"
Environment="PYTHONPATH=/home/capsule/capsulestory"
ExecStart=/home/capsule/capsulestory/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target